# crontab
42 2 * * * /usr/bin/certbot renew --quiet
